# containers/k3s/ssh-hostkey-regen.service
[Unit]
Description=Regenerate SSH host keys on first boot (per-VM uniqueness)
# Run before sshd would accept connections — sshd.service has its own keygen
# but we want a clean slate so the image-baked keys never travel further than
# the build host. ConditionPathExists makes this idempotent across reboots.
Before=sshd.service
ConditionPathExists=!/var/lib/ssh-host-keys-regenerated

[Service]
Type=oneshot
ExecStart=/usr/libexec/ssh-hostkey-regen.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
