# containers/k8s/k8s-init.service
[Unit]
Description=Initialize single-node kubeadm cluster on first boot
# Note: any After= on cloud-final.service or cloud-init.target creates a
# systemd ordering cycle (cloud-init's own units are inside multi-user.target's
# transaction; this service is WantedBy=multi-user.target). The cycle gets
# broken by systemd dropping k8s-init.service from the boot, so the cluster
# never initializes. Instead k8s-init.sh blocks on `cloud-init status --wait`
# at runtime when cloud-init is installed.
After=crio.service network-online.target
Wants=crio.service network-online.target
ConditionPathExists=!/var/lib/k8s-init.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/k8s-init.sh
RemainAfterExit=yes
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
