# containers/k8s/k8s-init.service
[Unit]
Description=Initialize single-node kubeadm cluster on first boot
# cloud-init.target is a no-op when cloud-init is not installed (the
# default build). When ENABLE_CLOUD_INIT=1 was passed at image build, this
# After= ensures user-data-injected SSH keys / runcmd / package-installs
# land before kubeadm init runs, so any operator-supplied prerequisites
# (extra packages, files written by write_files, etc.) are in place.
After=crio.service network-online.target cloud-final.service
Wants=crio.service network-online.target
ConditionPathExists=!/var/lib/k8s-init.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/k8s-init.sh
RemainAfterExit=yes
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
