# containers/shared/bootc-fetch-apply-updates.service.d/10-mark-upgrade.conf
# Drop-in for the upstream bootc-fetch-apply-updates.service. After a
# successful invocation (which stages the new deployment and then reboots),
# touch a marker file the next boot's health-check-rollback.timer reads to
# decide whether to run health checks + auto-rollback (#104).
#
# This is in /etc (not /usr) intentionally: we want the marker to be set as
# a side effect of the unit running, and dropping it under
# /etc/systemd/system/<unit>.service.d/ is the supported override path that
# doesn't require modifying the upstream unit.

[Service]
ExecStartPost=/usr/bin/touch /var/lib/bootc-just-upgraded
