# containers/k8s-worker/worker-init.service
[Unit]
Description=Join the Hummingbird Kubernetes cluster as a worker
# Any After= on cloud-final.service or cloud-init.target creates a systemd
# ordering cycle (those units are inside multi-user.target's transaction and
# this service is WantedBy=multi-user.target). worker-init.sh does NOT block
# on `cloud-init status --wait` either — that wait blocks on cloud-final,
# which is After=multi-user.target, which is itself blocked by this service
# (deadlock — PR #265). Hostname is read directly via `hostnamectl --static`
# because cloud-init's hostname module runs at the init stage, before
# multi-user.target activates.
After=crio.service network-online.target
Wants=crio.service network-online.target
ConditionPathExists=!/var/lib/worker-init.done
# Skip cleanly when an operator boots the bare template qcow2 (no token
# injected). spawn-workers.sh writes this file per-VM before virt-install.
ConditionPathExists=/etc/hummingbird/worker-join.env

[Service]
Type=oneshot
ExecStart=/usr/libexec/worker-init.sh
RemainAfterExit=yes
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
