# containers/k8s-worker/worker-init.service
[Unit]
Description=Join the Hummingbird Kubernetes cluster as a worker
# cloud-init.target is a no-op when cloud-init is not installed (the
# default build). When ENABLE_CLOUD_INIT=1 was passed at image build, this
# After= ensures user-data-injected SSH keys / runcmd / package-installs
# land before kubeadm join runs.
After=crio.service network-online.target cloud-final.service
Wants=crio.service network-online.target
ConditionPathExists=!/var/lib/worker-init.done
# Skip cleanly when an operator boots the bare template qcow2 (no token
# injected). spawn-workers.sh writes this file per-VM before virt-install.
ConditionPathExists=/etc/hummingbird/worker-join.env

[Service]
Type=oneshot
ExecStart=/usr/libexec/worker-init.sh
RemainAfterExit=yes
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
